# zfs-flush.timer
[Unit]
Description=Run ZFS flush script periodically

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
